# route.service.in
[Unit]
Description=route shadow
After=multi-user.target

[Service]
Type=oneshot
WorkingDirectory=@IPROUTE_GIT_DIR@
ExecStart=/usr/bin/srtunectl

[Install]
WantedBy=multi-user.target
